# scripts/vps-resilience/sost-conn-monitor.timer
[Unit]
Description=Snapshot SOST RPC port connection counts every hour
Requires=sost-conn-monitor.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h
AccuracySec=1min
Unit=sost-conn-monitor.service

[Install]
WantedBy=timers.target
