# scripts/vps-resilience/sost-conn-monitor.service
[Unit]
Description=SOST RPC port connection-count snapshot

[Service]
Type=oneshot
ExecStart=/opt/sost/scripts/monitor-connections.sh
